# recipes-support/rdmagent/files/apps-prerdm.service
[Unit]
Description= Service to ensure pre condition for rdm apps download
After=nvram.service

[Install]
WantedBy=multi-user.target
