# recipes-support/rdmagent/files/apps-rdm.service
[Unit]
Description= Apps Download Service
After=swupdate.service
ConditionPathExists=/tmp/.xconfssrdownloadurl

[Service]
SyslogIdentifier=apps-rdm
