# recipes-support/rdmagent/files/apps-rdm.service
[Unit]
Description= Apps Download Service
ConditionPathExists=/tmp/.xconfssrdownloadurl

[Service]
SyslogIdentifier=apps-rdm
ExecStart=/usr/bin/rdm
