# recipes-support/chrony/files/chrony-tracking.service
[Unit]
Description=Chrony telemetry data collection service

[Service]
Type=oneshot
ExecStart=/lib/rdk/chrony_tracking.sh
